# systemd/system/mdmonitor.service.d/00-syslog.conf
[Service]
Environment="MDADM_MONITOR_ARGS=--scan --syslog"

