# systemd/system/coreos-startup.target
[Unit]
Description=CoreOS essentials
Wants=issuegen.service
Wants=media.mount
Wants=usr-share-oem.mount

[Install]
WantedBy=basic.target
