# systemd/system/logrotate.timer
[Unit]
Description=Daily Log Rotation

[Timer]
OnCalendar=daily
AccuracySec=1h
Persistent=true
