# systemd/system/motdgen.path
[Unit]
Description=Watch for update engine configuration changes

[Path]
PathChanged=/etc/coreos/update.conf
