# systemd/system/motdgen.timer
[Unit]
Description=Generate /run/coreos/motd

[Timer]
OnUnitActiveSec=60
OnBootSec=0
