# systemd/system/motdgen.service
[Unit]
Description=Generate /run/coreos/motd
Before=systemd-user-sessions.service

[Service]
Type=oneshot
ExecStart=/usr/lib/coreos/motdgen
