# systemd/system/var.mount
[Unit]
Conflicts=umount.target
Before=local-fs.target umount.target

[Mount]
What=/media/state/overlays/var
Where=/var
Type=bind
Options=bind
