# systemd/system/sshd@.service.d/sshd-keygen.conf
[Unit]
Wants=sshd-keygen.service
After=sshd-keygen.service
