# systemd/system/systemd-fsck@dev-disk-by\x2dlabel-OEM.service.d/conditional.conf
[Unit]
ConditionVirtualization=!container
ConditionPathExists=!/usr/.noupdate
