# systemd/addon-config@.service
[Unit]
Description=Configure Addon /%I
ConditionPathExists=|/%I/authorized_keys

[Service]
Type=oneshot
ExecStart=/usr/lib/coreos/addon_config /%I
