# systemd/home.mount
[Unit]
Conflicts=umount.target
Before=local-fs.target umount.target

[Mount]
What=/media/state/overlays/home
Where=/home
Type=bind
Options=bind
