# systemd/srv.mount
[Unit]
Conflicts=umount.target
Before=local-fs.target umount.target

[Mount]
What=/media/state/overlays/srv
Where=/srv
Type=bind
Options=bind
